# deploy/sub/minikube/minikube.service
# Chinese article:
# https://www.ruanyifeng.com/blog/2016/03/systemd-tutorial-part-two.html

[Unit]
Description=minikube
After=docker.service

[Service]
ExecStart=sudo minikube start --vm-driver=none
ExecStop=minikube stop
Type=oneshot
RemainAfterExit=yes

[Install]
WantedBy=graphical.target